# deploy/systemd-user/inertia-ssr.service
[Unit]
Description=Streetcode Inertia SSR server
After=network.target
ConditionPathExists=/home/deployer/streetcode-laravel/current

[Service]
Type=simple
WorkingDirectory=/home/deployer/streetcode-laravel/current
ExecStart=/bin/bash -lc 'source ~/.nvm/nvm.sh 2>/dev/null; exec /usr/bin/php artisan inertia:start-ssr'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
